# root-armhf/etc/s6-overlay/s6-rc.d/init-config-kasm/run
#!/usr/bin/with-contenv bash

for i in {1..10}; do
  echo "Armhf/Arm32v7 is not supported"
done
