# root/etc/s6-overlay/s6-rc.d/svc-docker/run
#!/usr/bin/with-contenv bash

exec /usr/local/bin/dockerd-entrypoint.sh -l error --data-root /opt/docker
